# src/Pcntl/Conf/tail.ini
[program:tail]
command=/usr/bin/tail -f /home/xiakai/tttt
startsecs=3
autorestart=1